# crontab
*/15 * * * * ~/.local/bin/mailsync
0 * * * * calcurse-caldav
