# crontab
*/15 * * * * ~/.local/bin/mailsync
*/30 * * * * calcurse-caldav
